# dist/sys-oled.service
[Unit]
Description=SATA HAT System status OLED

[Service]
Type=exec
EnvironmentFile=-/etc/default/go-sata-hat
ExecStartPre=+/sbin/modprobe i2c-dev
ExecStart=/usr/bin/sys-oled

[Install]
WantedBy=default.target
